# docker/web/supervisord/conf.d/apache.conf
[program:apache]
command=apache2ctl -DFOREGROUND
autostart=true
autorestart=true

